# supervisor-app.ini
# NOTE: deploy to /etc/supervisor.d/
[program:app]
redirect_stderr = true
command = /app-docker-entrypoint.sh
